# deploy/brain/systemd/ai_generate_fun_facts.service
[Unit]
Description=Laravel Artisan Command for Generating Fun Facts
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/deploy/ai-things/current/api
ExecStart=/usr/bin/php /deploy/ai-things/current/api/artisan Ai:GenerateFunFacts --sleep=30
Restart=always
RestartSec=3
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=laravel-fun-facts

[Install]
WantedBy=multi-user.target
